# data/infographic-service.conf.in
description "Service to connect infographic data sources to infographic visualizations"

start on desktop-start
stop on desktop-end

emits untrusted-helper-start
emits untrusted-helper-end
emits untrusted-helper-type-end

respawn
respawn limit 5 10

exec @pkglibexecdir@/infographic-service

# Make sure if we go, we bring all our helpers with us
post-stop exec initctl emit untrusted-helper-type-end HELPER_TYPE=infographic-visualization
