# data/usermetricsservice.conf.in
description "Service to connect infographic data sources to infographic visualizations"

start on started dbus
stop on desktop-end

respawn
respawn limit 5 10

exec @CMAKE_INSTALL_FULL_LIBEXECDIR@/libusermetrics/usermetricsservice

